# dot_config/systemd/user/wallpaper-change.path
[Unit]
Description=Monitor wallpaper path file changes

[Path]
PathChanged=%h/.config/wallpaper-path

[Install]
WantedBy=graphical-session.target
